# root/etc/services.d/deluged/run
#!/usr/bin/with-contenv bash

umask 000

exec \
	s6-setuidgid abc /usr/bin/deluged -c /config \
	-d --loglevel=debug -l /config/deluged.log
